# homebrew.mxcl.vault.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>homebrew.mxcl.vault</string>

  <!-- Ensure launchd can find nc, vault, etc. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>
      # ---- BEGIN CONFIG ----
      PG_HOST="10.0.0.3"     
      PG_PORT="5432"         
      VAULT_BIN="/opt/homebrew/opt/vault/bin/vault"
      VAULT_HCL="/Users/isaac.christian/vault.d/vault.hcl"
      # ---- END CONFIG ----

      # Wait for network + Postgres to be reachable
      echo "Checking reachability for ${PG_HOST}:${PG_PORT} ..."
      until /usr/bin/nc -z "${PG_HOST}" "${PG_PORT}" >/dev/null 2>&amp;1; do
        echo "Waiting for remote Postgres at ${PG_HOST}:${PG_PORT} ..."
        sleep 5
      done
      echo "Remote Postgres reachable. Starting Vault."

      # Start Vault
      exec "${VAULT_BIN}" server -config "${VAULT_HCL}"
    </string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/isaac.christian/vault.d/logs/vault.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/isaac.christian/vault.d/logs/vault.err.log</string>

  <key>WorkingDirectory</key>
  <string>/opt/homebrew/var</string>
</dict>
</plist>